# package/harper/DEBIAN/postinst
#!/bin/bash

systemctl enable harper.service
systemctl start harper.service
